# quiz-tg.service
[Service]
ExecStart=/opt/quiz/venv/bin/python3 /opt/quiz/tg.py
Restart=always
Requires=quiz-data-upload.service
[Install]
WantedBy=multi-user.target